# supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf

[inet_http_server]
port=178.128.195.100:9001
username=admin
password=Af2f4dsASF

[program:openvidu]
command=/bin/bash -c "java -jar -Dopenvidu.secret=A21Qa#CnfqbCKQaHLqzT9LGruE -Dopenvidu.publicurl=https://simia.coupler.app:4443/ -Dserver.ssl.key-store=/usr/config/openvidu/keystore.jks -Dserver.ssl.key-store-password=Af2f4dsASF -Dserver.ssl.key-alias=simiaalias -Dopenvidu.recording=true -Dopenvidu.recording.path=/usr/storage/openvidu/records /usr/config/openvidu/openvidu-server-2.11.0.jar"
stdout_logfile=/var/log/openvidu.log
autostart=true
autorestart=true
redirect_stderr=true
priority=3
